# Docker/logrotate.cron
*/5 * * * * logrotate /etc/logrotate.d/syslog-ng --state /var/log/syslog-ng/logrotate.state