# Docker/logrotate.cron
*/5 * * * * logrotate /etc/logrotate.d/syslog-ng -d